# config/dnsovertls.conf
[program:dottcp]
command=/usr/bin/dottcp
stdout_logfile=/var/log/supervisor/dottcp.log

[program:dotudp]
command=/usr/bin/dotudp
stdout_logfile=/var/log/supervisor/dotudp.log